# default/etc/init.d/greetd
#!/sbin/openrc-run
name=$RC_SVCNAME
description="Greeter daemon"
supervisor="supervise-daemon"
command="greetd"